# src/daemon/com.later.muse.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.later.muse</string>

    <key>ProgramArguments</key>
    <array>
        <!-- 路径需根据实际环境调整 -->
        <string>/usr/local/bin/node</string>
        <string>/Users/xulater/Code/assistant-agent/muse/daemon/cerebellum.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/xulater/Code/assistant-agent</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/muse-cerebellum.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/muse-cerebellum.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
